# strategem-vps-minimal.service
[Unit]
Description=Strategem Core FastAPI Web Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/strategem
Environment="PATH=/opt/strategem/venv/bin"
Environment="OPENROUTER_API_KEY=YOUR_API_KEY_HERE"
Environment="LLM_MODEL=openai/gpt-4o-mini"
Environment="LLM_TEMPERATURE=0.2"
Environment="LLM_MAX_TOKENS=4000"
# Single worker for 1 core VPS, --reload disabled for production
ExecStart=/opt/strategem/venv/bin/uvicorn strategem.web.app:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
